# code/src/app/doc/app.conf
[program:app]
directory = /work/project
command = /work/project/app
autostart = true
startsecs = 5
user = root
redirect_stderr = true
stdout_logfile = /var/log/supervisor/app.log


